# .github/tmp/package/add/up/tool/luci-app-partexp/root/etc/uci-defaults/zzz_luci-app-partexp
#!/bin/sh

chmod +x /usr/bin/partexp /usr/libexec/rpcd/partexp

rm -f /tmp/luci-indexcache
exit 0
