# .github/tmp/package/add/up/wrtbwmonjs/luci-app-wrtbwmon/root/etc/uci-defaults/40-luci-wrtbwmon
#!/bin/sh

chmod +x /usr/libexec/rpcd/luci.wrtbwmon  2>/dev/null
rm -f /tmp/luci-indexcache

exit 0
